# installer/ada-pi-backend.service
[Unit]
Description=Ada-Pi Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ada-pi/backend
ExecStart=/usr/bin/python3 /opt/ada-pi/backend/main.py
Restart=always
RestartSec=5

# Ensure logging goes to systemd journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
